# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for peternicewicz.com"

start on net-device-up
stop on shutdown

respawn

setuid elspeth
chdir /home/elspeth/sites/peternicewicz.com/source

exec ../virtualenv/bin/gunicorn --bind unix:/tmp/peternicewicz.com.socket superlists.wsgi:application
